# scripts/auto-update-restart.timer
[Unit]
Description=Run auto update script daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
